# deploy/systemd/agp-backup.timer
[Unit]
Description=Run AGP backup daily

[Timer]
OnCalendar=*-*-* 02:30:00
Persistent=true
RandomizedDelaySec=900

[Install]
WantedBy=timers.target
